# examples/schedule/macos-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd job: re-index the photo library and deploy nightly.

  Install:
    cp examples/schedule/macos-launchd.plist \
       ~/Library/LaunchAgents/com.local.nas-gallery.plist
    # Edit the file: paste your actual paths into ProgramArguments,
    # WorkingDirectory, and the StandardOut/Err logs.
    launchctl load ~/Library/LaunchAgents/com.local.nas-gallery.plist

  Verify:
    launchctl list | grep nas-gallery
    tail -f /tmp/nas-gallery.log

  Stop:
    launchctl unload ~/Library/LaunchAgents/com.local.nas-gallery.plist

  Notes:
    - StartCalendarInterval below runs nightly at 04:30. If the Mac
      is asleep at that time, launchd will run the job when it next
      wakes (unless WakeFromSleep is set in the plist or pmset is
      configured to wake; both are out of scope here).
    - Use StartInterval (seconds) instead of StartCalendarInterval
      for "every N hours" semantics.
    - The wrapper script does build + deploy. Adjust paths to taste.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.nas-gallery</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/nas-gallery/examples/manual/Update_gallery.command</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/nas-gallery</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/nas-gallery.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/nas-gallery.err</string>

    <!-- Don't run multiple instances if a previous one is still going -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
